# deploy/zetrick-api.service
[Unit]
Description=Zetrick API
Documentation=https://github.com/dprandle/zetrick-api
After=network-online.target
Wants=network-online.target

[Service]
# Runs as root straight from the git checkout in root's home, so updates are a
# plain `git pull` followed by `systemctl restart zetrick-api` (the build step below
# recompiles dist/ on every (re)start).
WorkingDirectory=/root/zetrick-api

# Also cleans and rebuilds
ExecStart=/usr/bin/npm start

# Restart on crash; don't hammer on a hard config/build failure.
Restart=on-failure
RestartSec=10s
StartLimitIntervalSec=300
StartLimitBurst=10

# Give the in-process graceful shutdown 
# room to finish; its own 10s watchdog force-exits well before this.
TimeoutStopSec=20s

# Logs to the journal: journalctl -u qbtsync -f
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target